# udev/60-bluefruit-event.rules
KERNEL=="event*", SUBSYSTEM=="input", SUBSYSTEMS=="input", ATTRS{name}=="Bluefruit52_Joy Keyboard", ACTION=="add", RUN+="/bin/bluefruit-k2j.sh"
